# knowledge-base.service
[Unit]
Description=Knowledge Base API
After=network.target neo4j.service

[Service]
Type=simple
User=kbuser
Group=kbuser
WorkingDirectory=/opt/knowledge-base
ExecStart=/opt/knowledge-base/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
# JDK path: update to match your actual installation
Environment=JAVA_HOME=/opt/tools/jdk
Environment=PATH=/opt/knowledge-base/venv/bin:/opt/tools/jdk/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
